# port-tester-echo-udp.service
[Unit]
Description=port-tester-echo-udp

[Service]
ExecStart=/usr/bin/node /app/port-tester-echo-udp/index.js

[Install]
WantedBy=multi-user.target
